# _migrate/atl-sh/ansible/roles/monitoring/templates/node_exporter.service.j2
# {{ ansible_managed }}
[Service]
ExecStart=
ExecStart=/usr/bin/prometheus-node-exporter --web.listen-address={{ node_exporter_listen }}
