# _migrate/atl-sh/ansible/roles/services/templates/finger@.service.j2
[Unit]
Description=Finger Protocol Service
Documentation=man:efingerd(8)

[Service]
ExecStart=/usr/sbin/efingerd
StandardInput=socket
StandardOutput=socket
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
